# apps/hastuioctl/hastuioctl.service
[Unit]
Description=Hastuioctl — Home Assistant audio control daemon
Documentation=github:grim-toolbox/apps/hastuioctl
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env uv --no-project --with pyyaml --with paho-mqtt --with loguru --with pydantic run \
    %h/Projects/My/grim-toolbox/apps/hastuioctl/hastuioctl.py
WorkingDirectory=%h
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ReadWritePaths=%h/.config/hastuioctl
ProtectHome=true
NoNewPrivileges=true
# Allow access to D-Bus session bus and PipeWire/Pulse
BusPolicy=connection

[Install]
WantedBy=default.target
